# launchd/com.pokediscord.bridge.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.pokediscord.bridge</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd /Users/jake/Documents/Projects/PokeDiscord && POKE_AUTO_TUNNEL=false pnpm start</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/poke-discord-bridge.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/poke-discord-bridge.err</string>
</dict>
</plist>
